# config/templates/mcadmin_supervisor.conf
[program:mcadmin_websocket]
command=node start.js
process_name=%(program_name)s
numprocs=1
directory=/var/www/mcadmin/node
user=minecraft
autorestart=true