# config/templates/mooshroom_supervisor.conf
[program:mooshroom_{name}]
command={command}
process_name=%(program_name)s
numprocs=1
directory=/var/lib/mooshroom/current/node
user={user}
autorestart=true  