# config/templates/mooshroom_supervisor_logobserver.conf
[program:mooshroom_logobserver]
command=node logObserver.js
process_name=%(program_name)s
numprocs=1
directory={directory}
user={user}
autorestart=true  