# config/templates/mooshroom_supervisor.conf
[program:mooshroom_websocket]
command=node start.js
process_name=%(program_name)s
numprocs=1
directory={directory}
user={user}
autorestart=true  