# config/templates/supervisor.conf
[program:mc_{NAME}]
command=java -Xmx{MEMORY} -jar server.jar
process_name=%(program_name)s
numprocs=1
directory=~{user}/mc_{NAME}
user={user}

